# cloud_api/aws-deployment/systemd/ytdlp-api.service
[Unit]
Description=YouTube Size Extension Cloud API
Documentation=https://github.com/yourusername/ytdlp-api
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/ytdlp-api
User=ubuntu
Group=ubuntu

# Start command
ExecStart=/usr/bin/docker run -d \
    --name ytdlp-api \
    --restart unless-stopped \
    -p 3000:3000 \
    --env-file /opt/ytdlp-api/.env \
    -v /opt/ytdlp-api/logs:/app/logs \
    ytdlp-api:latest

# Stop command
ExecStop=/usr/bin/docker stop ytdlp-api
ExecStopPost=/usr/bin/docker rm ytdlp-api

# Restart policy
Restart=on-failure
RestartSec=10s

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
